# supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command="sleep inf"
