# supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command=/bin/sleep inf
