# supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command=/root/monit_daemon.sh