# supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command=/var/vcap/bosh/bin/monit -I
