# supervisord.conf
[supervisord]
nodaemon=true

[program:sleep]
command=/bin/sleep inf
