# supervisord.conf
[supervisord]
nodaemon=true

[program:sleep]
command=/root/monit_daemon.sh
